# whigi-giveaway/crontab
*/10 * * * * /usr/bin/find /var/www -name "wp-cron.php" -type f -exec php {} \;